# scripts/whispera-wiraid@.service
[Unit]
Description=Whispera WIRAID module: %i
After=network-online.target whispera.service
Wants=network-online.target
PartOf=whispera.service

[Service]
Type=simple
User=whispera
Group=whispera
WorkingDirectory=/var/lib/whispera/wiraid/%i
ExecStart=/var/lib/whispera/wiraid/%i/bin/%i --config /var/lib/whispera/wiraid/%i/config.json
Restart=on-failure
RestartSec=5

# Resource limits — adjust per module
MemoryMax=256M
MemoryHigh=192M
CPUQuota=50%
TasksMax=64
LimitNOFILE=4096
LimitNPROC=64

# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false

# Allow only the module's own working dir to be writable
ReadWritePaths=/var/lib/whispera/wiraid/%i

# Network: bind to loopback only by default
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @debug @cpu-emulation @obsolete @raw-io @reboot @swap
SystemCallErrorNumber=EPERM

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
